# gpu-server/scripts/memory-watchdog.service
[Unit]
Description=HeartCode Memory Watchdog Service
Documentation=https://github.com/heartcode
After=docker.service
Requires=docker.service
PartOf=docker.service

[Service]
Type=simple
ExecStart=/home/boss/heartcode/gpu-server/scripts/memory-watchdog.sh --daemon
Restart=always
RestartSec=10
Environment="WARN_THRESHOLD=80"
Environment="SHED_THRESHOLD=85"
Environment="EMERGENCY_THRESHOLD=90"
Environment="CRITICAL_THRESHOLD=95"
Environment="POLL_INTERVAL=5"
Environment="LOG_FILE=/var/log/memory-watchdog.log"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=memory-watchdog

# Must run as root to clear caches and restart containers
User=root
Group=docker

[Install]
WantedBy=multi-user.target
